# python-fan-pwm.service
[Unit]
Description=Regulates

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/pwm.py

[Install]
WantedBy=multi-user.target